# .config/supervisord.conf
[program:django]
command=python /srv/backend/app/manage.py runserver 0:8000

[program:front]
command=npm run start --prefix /srv/front

[program:nginx]
command=nginx